# debian/postrm
#!/bin/bash -e

if [ "$1" = "remove" ]; then
    update-rc.d -f dbcat remove
fi
