# debian/preinst
#!/bin/bash -e

if [ -f /etc/init.d/dbcat ]
then
    invoke-rc.d dbcat stop
fi
