# debian/prerm
#!/bin/bash -e


invoke-rc.d dbcat stop

rm -f /etc/nginx/sites-enabled/dbcat
invoke-rc.d nginx restart
